# pulseaudio/96-pulseaudio-usb.rules
ATTRS{idVendor}=="046d", ATTRS{idProduct}=="0ab7", RUN+="/usr/local/bin/relocate-audio-devices.sh"
